# conf.d/zs_deployment/supervisord/services-available/gf_sender.conf
[program:gf_sender]
command=python /home/pabb/ydcws/trunk/apps/gf_sender/server.py --port=40%(process_num)02d --mode=debug --logging=DEBUG
process_name=40%(process_num)02d
numprocs=1
numprocs_start=1
stopsignal=INT
stdout_logfile=/var/log/supervisor/%(program_name)s/access.log
stderr_logfile=/var/log/supervisor/%(program_name)s/error.log
