# conf.d/zs_deployment/supervisord/services-available/uweb.conf
[program:uweb]
command=python /home/pabb/acb_new/trunk/apps/uweb/server.py --port=80%(process_num)02d --mode=debug --logging=DEBUG
process_name=80%(process_num)02d
numprocs=1
numprocs_start=1
stopsignal=INT
stdout_logfile=/var/log/supervisor/%(program_name)s/access.log
stderr_logfile=/var/log/supervisor/%(program_name)s/error.log
