# conf.d/zs_deployment/supervisord/services-available/wspush.conf
[program:wspush]
directory=/home/pabb/admin
command=gulp
stdout_logfile=/var/log/supervisor/%(program_name)s/access.log
stderr_logfile=/var/log/supervisor/%(program_name)s/error.log

