# conf.d/zs_online/supervisord/services-available/sms.conf
[program:sms]
command=python /home/ydcws/acb/trunk/apps/sms/server.py --port=20%(process_num)02d --mode=debug --logging=DEBUG
process_name=20%(process_num)02d
numprocs=1
numprocs_start=1
stopsignal=INT
stdout_logfile=/var/log/supervisor/%(program_name)s/access.log
stderr_logfile=/var/log/supervisor/%(program_name)s/error.log
