# conf.d/zs_online/supervisord/services-available/wspush.conf
[program:wspush]
directory=/home/ydcws/admin
command=gulp no-watch
stdout_logfile=/var/log/supervisor/%(program_name)s/access.log
stderr_logfile=/var/log/supervisor/%(program_name)s/error.log

